# recorder/systemd/airpoc-recorder.service
[Unit]
Description=AIRPOC record & replay daemon (:8093)
After=network.target data.mount
# Hard dependency on the recording volume: without this a failed data.mount let the
# daemon start anyway and "record" into an empty root-fs path — silent data loss.
RequiresMountsFor=/data/recordings

[Service]
Type=simple
User=asaftg
ExecStart=/home/asaftg/AIRPOC/recorder/src/airpoc_recorder -p 8093 -r /data/recordings
Restart=always
RestartSec=2
Nice=5

[Install]
WantedBy=multi-user.target
